# resources/supervisor/beta/supervisord.conf
# PUT THIS FILE IN /etc/supervisord.conf

[supervisord]
logfile=/data1/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:es]
command=/usr/local/elasticsearch/bin/elasticsearch
directory=/usr/local/elasticsearch
autostart=true
autorestart=true
startretries=1000000
stopsignal=INT
stopwaitsecs=60
stderr_logfile=/usr/local/elasticsearch/logs/es.error.log
stdout_logfile=/usr/local/elasticsearch/logs/es.log
user=klahnakoski
environment=ES_MIN_MEM="5g";ES_MAX_MEM="5g"

[program:etl]
command=python testlog_etl/etl.py --settings=resources/settings/beta/etl.json
directory=/home/ubuntu/TestLog-ETL
autostart=true
autorestart=true
startretries=1000000
stopsignal=INT
stopwaitsecs=10
stderr_logfile=/home/ubuntu/TestLog-ETL/results/logs/etl.error.log
stdout_logfile=/home/ubuntu/TestLog-ETL/results/logs/etl.log
user=klahnakoski
environment=PYTHONPATH="."

[program:monitor_es]
command=python testlog_etl/monitor/es.py --settings=resources/settings/beta/monitor_es.json
directory=/home/ubuntu/TestLog-ETL
autostart=true
startintervalsecs=120
autorestart=true
startretries=1000000
stopsignal=INT
stopwaitsecs=10
stderr_logfile=/home/ubuntu/TestLog-ETL/results/logs/monitor_es.error.log
stdout_logfile=/home/ubuntu/TestLog-ETL/results/logs/monitor_es.log
user=ubuntu
environment=PYTHONPATH='.'
