# resources/cron/manager.cron

* * * * * ~\TestLog-ETL\resources\scripts\run_buildbot_json_logs.sh
